# backend/data/mitre_techniques.json
[
  {"id": "T1595", "name": "Active Scanning", "tactic": "Reconnaissance", "tactic_id": "TA0043", "description": "Adversaries may execute active reconnaissance scans to gather information that can be used during targeting.", "detection": "Monitor for suspicious network traffic patterns indicating scanning activity.", "platforms": ["PRE"], "url": "https://attack.mitre.org/techniques/T1595/"},
  {"id": "T1590", "name": "Gather Victim Network Information", "tactic": "Reconnaissance", "tactic_id": "TA0043", "description": "Adversaries may gather information about the victim's networks that can be used during targeting.", "detection": "Monitor for lookups of domain registration data, DNS records, and network topology.", "platforms": ["PRE"], "url": "https://attack.mitre.org/techniques/T1590/"},
  {"id": "T1589", "name": "Gather Victim Identity Information", "tactic": "Reconnaissance", "tactic_id": "TA0043", "description": "Adversaries may gather information about the victim's identity that can be used during targeting (emails, credentials, names).", "detection": "Monitor for credential stuffing attempts and data broker queries.", "platforms": ["PRE"], "url": "https://attack.mitre.org/techniques/T1589/"},
  {"id": "T1596", "name": "Search Open Technical Databases", "tactic": "Reconnaissance", "tactic_id": "TA0043", "description": "Adversaries may search freely available technical databases for information about victims (Shodan, Censys, CT logs).", "detection": "Difficult to detect — passive reconnaissance against third-party databases.", "platforms": ["PRE"], "url": "https://attack.mitre.org/techniques/T1596/"},
  {"id": "T1592", "name": "Gather Victim Host Information", "tactic": "Reconnaissance", "tactic_id": "TA0043", "description": "Adversaries may gather information about the victim's hosts (OS, hardware, software, patches).", "detection": "Monitor for banner grabbing and service probing activity.", "platforms": ["PRE"], "url": "https://attack.mitre.org/techniques/T1592/"},
  {"id": "T1046", "name": "Network Service Discovery", "tactic": "Discovery", "tactic_id": "TA0007", "description": "Adversaries may attempt to get a listing of services running on remote hosts and local network infrastructure devices.", "detection": "Monitor for unusual port scanning activity, sequential connection attempts.", "platforms": ["Windows", "Linux", "macOS"], "url": "https://attack.mitre.org/techniques/T1046/"},
  {"id": "T1082", "name": "System Information Discovery", "tactic": "Discovery", "tactic_id": "TA0007", "description": "Adversaries may attempt to get detailed information about the operating system and hardware.", "detection": "Monitor for execution of system information commands (systeminfo, uname, cat /etc/os-release).", "platforms": ["Windows", "Linux", "macOS"], "url": "https://attack.mitre.org/techniques/T1082/"},
  {"id": "T1087", "name": "Account Discovery", "tactic": "Discovery", "tactic_id": "TA0007", "description": "Adversaries may attempt to get a listing of accounts on a system or within an environment.", "detection": "Monitor for commands like net user, Get-ADUser, cat /etc/passwd.", "platforms": ["Windows", "Linux"], "url": "https://attack.mitre.org/techniques/T1087/"},
  {"id": "T1135", "name": "Network Share Discovery", "tactic": "Discovery", "tactic_id": "TA0007", "description": "Adversaries may look for folders and drives shared on remote systems to identify sources of information.", "detection": "Monitor for net share, smbclient, and similar enumeration commands.", "platforms": ["Windows", "Linux"], "url": "https://attack.mitre.org/techniques/T1135/"},
  {"id": "T1190", "name": "Exploit Public-Facing Application", "tactic": "Initial Access", "tactic_id": "TA0001", "description": "Adversaries may attempt to exploit a weakness in an Internet-facing host or program.", "detection": "Monitor application logs for unusual requests, WAF alerts, error spikes.", "platforms": ["Windows", "Linux"], "url": "https://attack.mitre.org/techniques/T1190/"},
  {"id": "T1110", "name": "Brute Force", "tactic": "Credential Access", "tactic_id": "TA0006", "description": "Adversaries may use brute force techniques to gain access to accounts when passwords are unknown.", "detection": "Monitor for multiple failed login attempts, account lockouts.", "platforms": ["Windows", "Linux"], "url": "https://attack.mitre.org/techniques/T1110/"},
  {"id": "T1003", "name": "OS Credential Dumping", "tactic": "Credential Access", "tactic_id": "TA0006", "description": "Adversaries may attempt to dump credentials to obtain account login and credential material (hashes, tickets).", "detection": "Monitor for LSASS access, SAM database reads, DCSync traffic.", "platforms": ["Windows", "Linux"], "url": "https://attack.mitre.org/techniques/T1003/"},
  {"id": "T1078", "name": "Valid Accounts", "tactic": "Persistence", "tactic_id": "TA0003", "description": "Adversaries may obtain and abuse credentials of existing accounts as a means of gaining access.", "detection": "Monitor for unusual login times, locations, and impossible travel.", "platforms": ["Windows", "Linux"], "url": "https://attack.mitre.org/techniques/T1078/"},
  {"id": "T1059", "name": "Command and Scripting Interpreter", "tactic": "Execution", "tactic_id": "TA0002", "description": "Adversaries may abuse command and script interpreters to execute commands, scripts, or binaries.", "detection": "Monitor for unusual PowerShell, Bash, Python execution patterns.", "platforms": ["Windows", "Linux", "macOS"], "url": "https://attack.mitre.org/techniques/T1059/"},
  {"id": "T1053", "name": "Scheduled Task/Job", "tactic": "Persistence", "tactic_id": "TA0003", "description": "Adversaries may abuse task scheduling functionality to facilitate initial or recurring execution of malicious code.", "detection": "Monitor for new scheduled tasks, cron jobs, at jobs.", "platforms": ["Windows", "Linux"], "url": "https://attack.mitre.org/techniques/T1053/"},
  {"id": "T1548", "name": "Abuse Elevation Control Mechanism", "tactic": "Privilege Escalation", "tactic_id": "TA0004", "description": "Adversaries may circumvent mechanisms designed to control elevated privileges (UAC bypass, sudo abuse, SUID).", "detection": "Monitor for UAC bypass attempts, unusual sudo usage, SUID binary execution.", "platforms": ["Windows", "Linux"], "url": "https://attack.mitre.org/techniques/T1548/"},
  {"id": "T1068", "name": "Exploitation for Privilege Escalation", "tactic": "Privilege Escalation", "tactic_id": "TA0004", "description": "Adversaries may exploit software vulnerabilities to escalate privileges.", "detection": "Monitor for kernel exploit indicators, unusual process crashes, privilege changes.", "platforms": ["Windows", "Linux"], "url": "https://attack.mitre.org/techniques/T1068/"},
  {"id": "T1021", "name": "Remote Services", "tactic": "Lateral Movement", "tactic_id": "TA0008", "description": "Adversaries may use valid accounts to log into remote services (RDP, SSH, SMB, WinRM).", "detection": "Monitor for unusual remote login patterns, new source IPs.", "platforms": ["Windows", "Linux"], "url": "https://attack.mitre.org/techniques/T1021/"},
  {"id": "T1558", "name": "Steal or Forge Kerberos Tickets", "tactic": "Credential Access", "tactic_id": "TA0006", "description": "Adversaries may attempt to subvert Kerberos authentication (Kerberoasting, Golden/Silver tickets).", "detection": "Monitor for unusual TGS requests, RC4 encryption in Kerberos traffic.", "platforms": ["Windows"], "url": "https://attack.mitre.org/techniques/T1558/"},
  {"id": "T1205", "name": "Traffic Signaling", "tactic": "Defense Evasion", "tactic_id": "TA0005", "description": "Adversaries may use traffic signaling to hide open ports or trigger special actions on a system.", "detection": "Monitor for unusual packet patterns preceding port openings.", "platforms": ["Windows", "Linux"], "url": "https://attack.mitre.org/techniques/T1205/"}
]
